# crontab
*/5 * * * * docker run --rm your-container /root/update_table.py